# bot/deploy/stateof1c.service
# Юнит systemd для бота опроса (bot/README.md, раздел «Развертывание»).
# Установка: скопировать в /etc/systemd/system/stateof1c.service,
# токен и соль — в /etc/stateof1c.env (BOT_TOKEN=..., BOT_SALT=...).
[Unit]
Description=Бот опроса @stateOf1c_bot (landscape1c.ru)
After=network-online.target
Wants=network-online.target

[Service]
User=stateof1c
WorkingDirectory=/opt/landscape1c
ExecStart=/usr/bin/node bot/bot.js
EnvironmentFile=/etc/stateof1c.env
Restart=always
RestartSec=3
# SIGTERM при остановке — бот досохраняет state.json (обработчик в lib/store.js)
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
